# packaging/freebsd/pkg-plist
bin/dmarc-report-milter
etc/rc.d/milter-dmarc-report
@dir /var/run/milter-dmarc-report
